# .config/systemd/user/random-wallpaper.service
[Unit]
Description=set a random wallpaper

[Service]
Type=oneshot
ExecStart=%h/.config/i3/random_wallpaper.py -o %h/.cache/wallpaper.png

